# linux/01-evic.rules
SUBSYSTEM=="block", ACTION=="add", ATTRS{idVendor}=="0603", ATTRS{idProduct}=="b5d1", RUN+="/bin/sh /home/stefanb/bin/myvapors-parser/linux/evic_connected $env{ACTION} $env{DEVNAME}"
